# config/supervisord.d/exim.ini
[program:exim]
command=/usr/sbin/exim -bd -q5m -oX 127.0.0.1:1025
user=app
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0